# templates/swayidle.service.j2
[Unit]
Description = Swayidle Service

[Service]
Type=simple
ExecStart=/home/dhoessl/.config/sway/swayidle.sh
Restart=on-failure

[Install]
WantedBy=default.target
